# install_vmess_ipv6.sh
#!/bin/sh

set -e

XRAY_VERSION="1.8.24"
XRAY_URL="https://github.com/XTLS/Xray-core/releases/download/v${XRAY_VERSION}/Xray-linux-64.zip"

# ä½ çš„ IPv6 åœ°å€
SERVER_IP6="2001:41d0:303:3e79:be24:11ff:feff:ac09"

WS_PORT=8080
WS_PATH="/ws"
UUID=$(cat /proc/sys/kernel/random/uuid)

echo "ðŸ“¦ å®‰è£…ä¾èµ–..."
apk update
apk add --no-cache curl unzip bash openrc

echo "ðŸ“¥ ä¸‹è½½ Xray..."
mkdir -p /usr/local/xray
cd /usr/local/xray
curl -L -o xray.zip $XRAY_URL
unzip xray.zip
chmod +x xray

echo "âš™ï¸ å†™å…¥ Xray é…ç½®..."
mkdir -p /etc/xray

cat > /etc/xray/config.json <<EOF
{
  "inbounds":[
    {
      "listen": "::",
      "port": ${WS_PORT},
      "protocol": "vmess",
      "settings":{
        "clients":[{"id":"${UUID}"}]
      },
      "streamSettings":{
        "network":"ws",
        "wsSettings":{
          "path":"${WS_PATH}"
        }
      }
    }
  ],
  "outbounds":[
    {
      "protocol":"freedom",
      "settings":{}
    }
  ]
}
EOF

echo "ðŸ›  åˆ›å»º Alpine OpenRC æœåŠ¡..."
cat > /etc/init.d/xray <<EOF
#!/sbin/openrc-run

name="xray"
command="/usr/local/xray/xray"
command_args="run -c /etc/xray/config.json"
command_background=true
pidfile="/run/\${RC_SVCNAME}.pid"

depend() {
  need net
}
EOF

chmod +x /etc/init.d/xray
rc-update add xray default

echo "ðŸš€ å¯åŠ¨ Xray æœåŠ¡..."
rc-service xray restart

echo ""
echo "ðŸŽ‰ å®‰è£…æˆåŠŸï¼ï¼ä¸‹é¢æ˜¯ä½ çš„å®¢æˆ·ç«¯ä¿¡æ¯ï¼š"
echo "â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”"
echo "æœåŠ¡å™¨åœ°å€ : [${SERVER_IP6}]"
echo "ç«¯å£       : ${WS_PORT}"
echo "UUID       : ${UUID}"
echo "åè®®       : VMess"
echo "ä¼ è¾“æ–¹å¼   : WebSocket"
echo "WSè·¯å¾„     : ${WS_PATH}"
echo "TLS        : å…³é—­"
echo "â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”"
echo ""
echo "ðŸ“Œ æ³¨æ„ï¼š"
echo "åœ¨å®¢æˆ·ç«¯å¡«å†™ IPv6 æ—¶ï¼Œåœ°å€è¦åŠ ä¸­æ‹¬å·ï¼š"
echo "ä¾‹å¦‚ï¼š [${SERVER_IP6}]"
